# service/com.claude-note.worker.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.claude-note.worker</string>

    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>{{BIN_PATH}}/claude-note</string>
        <string>worker</string>
    </array>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it dies -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>{{VAULT_PATH}}/.claude-note/logs/worker-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{VAULT_PATH}}/.claude-note/logs/worker-stderr.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>{{HOME}}/.local/share/claude-note/src</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:{{HOME}}/.local/bin</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>{{HOME}}</string>

    <!-- Process limits -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>

    <!-- Throttle restarts (10 second minimum between respawns) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
